# docker/supervisord.conf
; ============================================================
; Supervisor Configuration for URL Shortener Backend
; ============================================================
; This configuration manages the PHP-FPM process
; ============================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.conf
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=www-data

; Emergency restart if PHP-FPM crashes
exitcodes=0,78,143
stopsignal=QUIT
stopwaitsecs=10
